# run-timelapse-allsky.timer
[Unit]
Description=Periodically generate and upload AllSky timelapse

[Timer]
OnCalendar=*-*-* 12:30:00

[Install]
WantedBy=multi-user.target
